# build_deb/scripts/prerm
#!/bin/bash

/etc/init.d/rocksserver stop
update-rc.d rocksserver remove

